# config/systemd/user/mail-sync-notify.service
[Unit]
Description=Sync QQ mail and send desktop notifications

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mail-sync-notify
