# config/autostart/com.abdownloadmanager.desktop
[Desktop Entry]
Type=Application
Name=AB Download Manager
Exec="/usr/bin/ABDownloadManager" --background
Terminal=false
NoDisplay=true
